# htracker-server/htracker.service
[Unit]
Description=Htracker HTTP Server

[Service]
ExecStart=/bin/bash /home/grant/htracker

[Install]
WantedBy=multi-user.target